# manager-centos6.5/components/restservice/config/upstart/cloudify-restservice.conf
description "restservice"

start on runlevel [2345]
stop on runlevel [!2345]

env PIDFile=/var/run/gunicorn.pid
env REST_SERVICE_CONFIG_PATH="/opt/manager/cloudify-rest.conf"
env MANAGER_REST_CONFIG_PATH="/opt/manager/cloudify-rest.conf"
env MANAGER_REST_SECURITY_CONFIG_PATH="/opt/manager/rest-security.conf"
env MANAGER_VIRTUALENV="/opt/manager/env"
# REST Service (and Gunicorn) Log Path
env REST_SERVICE_LOG_PATH="/var/log/cloudify/rest"
# REST Service port
env REST_SERVICE_PORT=8100
env REST_SERVICE_HOME="/opt/manager"
env REST_SERVICE_VIRTUALENV="/opt/manager/env"
env MANAGER_VIRTUALENV="/opt/manager/env"

respawn
respawn limit 10 5

exec /opt/manager/env/bin/gunicorn \
          --pid /var/run/gunicorn.pid \
           -w $(($(nproc)*2+1)) \
           -b 0.0.0.0:${REST_SERVICE_PORT} \
          --timeout 300 manager_rest.server:app \
          --log-file /var/log/cloudify/rest/gunicorn.log \
          --access-logfile /var/log/cloudify/rest/gunicorn-access.log
